# dummy_java/supervisord.conf
[supervisord]
nodaemon=true

[program:java-app]
command=java -jar /app/app.jar

[program:nginx]
command=nginx -g 'daemon off;' 